# .config/autostart/nitrogen.desktop
[Desktop Entry]
Version=1.0
Name=Nitrogen
Exec=nitrogen --restore
StartupNotify=false
Type=Application
